# crontab
0 5 * * * /code/manage.py enqueue_tasks --task run_routines