# crontab
0 5 * * * /code/manage.py enqueue_tasks --task run_routines
0 5 * * * /code/manage.py enqueue_tasks --task recalculate_contact_strengths